# config/supervisor-app.conf
[supervisord]
[program:app-uwsgi]
command = /venv/bin/uwsgi --ini /var/html/config/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx
